# driver-bot-simple.service
[Unit]
Description=Driver Bot Simple
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/driver-bot
Environment=PATH=/opt/driver-bot/venv/bin
ExecStart=/opt/driver-bot/venv/bin/python run_simple_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 